# config/includes.chroot/etc/profile.d/custom_path.sh
# Add /usr/local/bin to the beginning of the PATH
export PATH=/usr/local/bin:$PATH

